# rpi/gmusic-uploader/container/cron/root.crontab
*/5 * * * * /usr/bin/flock -n /tmp/fcj.lockfile -c "/scripts/sync.sh /credentials $UPLOADER_ID"
 
